# build/agent/service-discover.service
# SPDX-FileCopyrightText: 2021-2022, 2024-2025 Zextras <https://www.zextras.com>
#
# SPDX-License-Identifier: AGPL-3.0-only

[Unit]
Description=Zextras Service Discover Agent
Documentation=https://docs.zextras.com/
Requires=network-online.target
After=network-online.target
PartOf=service-discover.target

[Service]
Type=notify
ExecStart=/usr/bin/service-discoverd agent
Restart=on-failure
RestartSec=15s
KillMode=process
LimitNOFILE=65536
ExecReload=/usr/bin/consul reload
TimeoutSec=120s
TimeoutStopSec=120s

# Hardening
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=/var/lib/service-discover /etc/zextras/service-discover /run/carbonio /tmp
PrivateDevices=yes
ProtectHome=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
RestrictRealtime=yes
RestrictNamespaces=yes
RestrictSUIDSGID=yes
LockPersonality=yes
ProtectHostname=yes
ProtectClock=yes
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX AF_NETLINK

[Install]
WantedBy=service-discover.target
