# build/agent/service-discover.service
# SPDX-FileCopyrightText: 2021-2022, 2024-2025 Zextras <https://www.zextras.com>
#
# SPDX-License-Identifier: AGPL-3.0-only

[Unit]
Description=Service discover agent for Zextras
Documentation=https://docs.zextras.com/
Requires=network-online.target
After=network-online.target
PartOf=carbonio.target

[Service]
Type=notify
ExecStart=/usr/bin/service-discoverd agent
Restart=on-failure
RestartSec=15
KillMode=process
LimitNOFILE=65536
ExecReload=/usr/bin/consul reload
TimeoutSec=120
TimeoutStopSec=120

[Install]
WantedBy=multi-user.target
